# daemon-scripts/upstart/aria-server.conf
description "Aria Server"
author  "Abil N George <mail@abilng.in>"

start on runlevel [234]
stop on runlevel [0156]

chdir /usr/share/aria-server
exec /usr/share/aria-server/main.py
respawn
